# .config/systemd/user/notification-sound.service
[Unit]
Description=Notification Sound Service
After=graphical-session.target
Requires=dbus.service

[Service]
Type=simple
ExecStart=%h/.config/swaync/notification-sound.sh
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=PULSE_SERVER=unix:%t/pulse/native

[Install]
WantedBy=default.target