# .config/systemd/user/sunshine.service.d/override.conf
[Service]
ExecStartPre=
ExecStartPre=/bin/bash -c 'until hyprctl version 2>/dev/null; do sleep 1; done; sleep 5'
ExecStart=
ExecStart=/bin/bash -c 'hyprctl dispatch exec /usr/bin/sunshine'
